# SENAO/package/fingerprint/files/etc/init.d/zfingerprint_v4.init
#!/bin/sh /etc/rc.common
. /etc/functions.sh
START=99

start() {
	/usr/shc/zfingerprint_v4 start
}

stop() {
	/usr/shc/zfingerprint_v4 stop
}

