# SENAO/feeds/packages/net/ifplugd/sn_files/ifplugd.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=50
STOP=50


start() {
	/usr/shc/ifplugd start
}

stop() {
	/usr/shc/ifplugd stop
}
