# SENAO/vendors/EnGenius/EAP1300/zzSENAO_VENDOR/files/etc/init.d/boot
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=10
STOP=98

start() {
	/usr/shc/boot start
}

stop() {
	/usr/shc/boot stop
}
