# SENAO/vendors/EnGenius/EAP1300/zzSENAO_VENDOR/files/etc/init.d/factory_setting
#!/bin/sh /etc/rc.common

START=19

start() {
	/usr/shc/factory_setting start
}

stop() {
	/usr/shc/factory_setting stop
}
