# SENAO/vendors/EnGenius/EAP1300/zzSENAO_VENDOR/files/etc/init.d/mount-nand-ubifs
#!/bin/sh /etc/rc.common

START=99

start() {
	/usr/shc/mount-nand-ubifs start
}
