# SENAO/vendors/EnGenius/EAP1300/zzSENAO_VENDOR/files/etc/init.d/telnet
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=50

start() {
	/usr/shc/telnet start
}

stop() {
	/usr/shc/telnet stop
}
