# SENAO/package/senao-openapi-server/files/etc/hotplug.d/net/99-app-reload-avahi
if [ "${INTERFACE#app}" != "${INTERFACE}" ]
then
	if [ "$ACTION" = "add" ]
	then
		/etc/init.d/avahi-daemon restart
	fi
fi
