# files/etc/hotplug.d/net/50-mesh-led
if [ "$INTERFACE" = "bat0" ]; then
	/etc/init.d/led reload 
fi
