# files/etc/hotplug.d/ifplugd/05-ifplugd_lan_led
#!/bin/sh
. /lib/functions.sh
. /lib/functions/network.sh


	luci-reload led
