# files/etc/hotplug.d/ifplugd/20-mgmt-release
#!/bin/sh
TEMP_PRE_IP_FILE="/tmp/pre_ip"

[ $ACTION == "down" ] && {
    rm $TEMP_PRE_IP_FILE
}
