# SENAO/package/senao-syseye/files/syseye.init
#!/bin/sh /etc/rc.common

. /lib/functions.sh

START=13  ## after snlogd
STOP=97  ## before syslogd

boot() {
    start
}

start() {
    /usr/sbin/syseye -d
}

stop() {
    killall -15 syseye
}

reload() {
    stop
    start
}
