# files/etc/tp.d/T99rmgmt
#!/bin/sh /etc/rc.common

TP=99

start() {
        /sbin/rmgmt -p -d -u &
}
