# files/etc/uci-defaults/luci-i18n-german
#!/bin/sh
uci batch <<-EOF
	set luci.languages.de=German
        commit luci
EOF

