# files/etc/uci-defaults/luci-i18n-french
#!/bin/sh
uci batch <<-EOF
	set luci.languages.fr='French'
        commit luci
EOF

