# files/etc/uci-defaults/luci-theme-web_ECB_FANCY
#!/bin/sh
uci batch <<-EOF
	set luci.themes.web_ECB_FANCY=/luci-static/web_ECB_FANCY
	commit luci
EOF
exit 0
