# files/etc/uci-defaults/remove_wifiprofile
#!/bin/sh

# Remove wifiprofile let qcawifi.sh create again
rm /etc/config/wifiprofile
sync

exit 0
