# files/etc/init.d/led
#!/bin/sh /etc/rc.common
# (C) 2008 openwrt.org

START=21

start() {
	/usr/shc/led start
}
