# files/etc/init.d/locator
#!/bin/sh /etc/rc.common

START=99
#STOP=10

start() {
	/usr/shc/locator start
}

stop() {
	/usr/shc/locator stop
}

