# SENAO/package/senao-syskey/files/usr/sbin/iface_mark_boot.sh
#!/bin/sh /etc/rc.common

START=21

boot()
{
    /usr/sbin/iface_mark.sh
}
