# files/usr/lib/opkg/info/kmod-qca-rfs.list
/lib/modules/4.4.60/qrfs.ko
/etc/init.d/qrfs
/etc/modules.d/29-qca-rfs
