# files/usr/lib/opkg/info/kmod-qca-wifi-unified-profile.list
/etc/assoc_notify.d/README
/lib/wifi/apLOption.csv
/usr/sbin/thermaltool
/usr/sbin/exttool
/usr/sbin/athadhoc
/lib/wifi/RegularDomain.sh
/sbin/wifi_reload_log.sh
/lib/wifi/sn_wifi_version/version-apsteer
/sbin/BackgroundScan.sh
/lib/modules/4.4.60/ath_rate_atheros.ko
/lib/wifi/sn_wifi_version/version-senaonasport
/lib/modules/4.4.60/ath_pktlog.ko
/lib/wifi/sn_wifi_version/version-reloadboost
/usr/sbin/apstats
/lib/modules/4.4.60/qdf.ko
/sbin/sn_wifi_version
/lib/wifi/sn_wifi_version/version-wpa3
/lib/wifi/sn_wifi_version/version-chanspec
/usr/sbin/spectraltool
/lib/modules/4.4.60/qca_ol.ko
/lib/modules/4.4.60/ath_dfs.ko
/sbin/updateATFConfig.sh
/lib/wifi/nawds.sh
/sbin/chkEth2GW.sh
/etc/config/sn_reconn
/sbin/ScanRadio.sh
/sbin/Fastscan24G.sh
/usr/sbin/pktlogdump
/sbin/kicksta.sh
/usr/sbin/80211stats
/lib/wifi/sn_wifi_version/version-bandsteering
/sbin/layer_reload.sh
/lib/wifi/ignoreOption.csv
/lib/wifi/sn_wifi_version/version-hotspot2
/sbin/getWifiFirstIF
/sbin/boost_reload.sh
/lib/wifi/sn_wifi_version/version-meshnawds
/lib/modules/4.4.60/ath_spectral.ko
/lib/wifi/sn_wifi_version/version-scanningradio
/lib/modules/4.4.60/smart_antenna.ko
/lib/modules/4.4.60/ath_dev.ko
/usr/sbin/radartool
/sbin/sn_reconn
/lib/wifi/qcawifi_countrycode.txt
/lib/wifi/cmdLOption.csv
/sbin/11kPreScan.sh
/usr/sbin/wifitool
/usr/lib/libtlvparser.so
/lib/wifi/sn_wifi_version/version-backgroundscan
/usr/lib/libtlvencoder.so
/sbin/Fastscan5G.sh
/usr/sbin/athtestcmd
/sbin/checkOption.sh
/lib/wifi/sn_wifi_version/version-senao11k
/usr/sbin/tx99tool
/lib/wifi/sn_wifi_version/version-senao11r
/lib/modules/4.4.60/asf.ko
/sbin/proxyarp_clear.sh
/usr/sbin/athstats
/lib/wifi/vlan_passthrough.sh
/lib/modules/4.4.60/ath_hal.ko
/usr/lib/libtlvtemplate.so
/etc/init.d/sn_reconn
/usr/lib/libtlvcmdrsp.so
/sbin/sn_bandsteering
/sbin/assoc_notify.sh
/usr/lib/libtlvutil.so
/usr/sbin/wlanconfig
/lib/wifi/radioLOption.csv
/lib/modules/4.4.60/hst_tx99.ko
/lib/wifi/wifi_profile
/lib/modules/4.4.60/umac.ko
/sbin/setChannelConfig.sh
/lib/modules/4.4.60/qca_da.ko
/lib/wifi/vapLOption.csv
/usr/sbin/athstatsclr
/usr/sbin/ssidsteering
/usr/sbin/wps_enhc
/usr/sbin/pktlogconf
/lib/wifi/RegularDomain.csv
/lib/wifi/sn_wifi_version/version-wificommon
/etc/modules.d/33-qca-wifi-unified-profile
/lib/wifi/wifi_funcs.sh
/lib/modules/4.4.60/mem_manager.ko
/lib/wifi/qcawifi.sh
