# files/usr/lib/opkg/info/kmod-qca-edma.list
/etc/modules.d/45-qca-edma
/etc/init.d/qca-edma
/lib/modules/4.4.60/essedma.ko
