# files/usr/lib/opkg/info/locator.list
/usr/shc/locator
/etc/uci-defaults/locator
/etc/init.d/locator
/usr/sbin/discover
/usr/sbin/locator
