# files/usr/lib/opkg/info/openlldp.list
/etc/uci-defaults/openlldpd
/etc/init.d/openlldpd
/usr/bin/lldpd
