# files/usr/lib/opkg/info/qca-hostap.list
/lib/wifi/hostapd_profile
/etc/init.d/qca-hostapd
/usr/sbin/hostapd
/lib/wifi/hostapd.sh
/usr/sbin/wpa-psk-check.sh
