# files/usr/lib/opkg/info/qca-thermald-10.4.list
/usr/sbin/thermald
/etc/config/thermal
/etc/thermal/ipq-thermald-8069.conf
/etc/init.d/thermal
/etc/thermal/ipq-thermald-8064.conf
/etc/thermal/ipq-thermald-8066.conf
