# files/usr/lib/opkg/info/rmgmt.list
/etc/modelname
/etc/init.d/rmgmt
/sbin/rmgmt
