# files/usr/lib/opkg/info/senao-l2-isolation.list
/etc/uci-defaults/l2_isolation
/etc/init.d/l2_isolation
/usr/sbin/dhtest_timer
/etc/config/l2_isolation
