# files/usr/lib/opkg/info/senao-guest-syncd.list
/usr/sbin/guest_syncli
/etc/uci-defaults/guestsyncd
/etc/init.d/guestsyncd
/etc/config/guestsyncd
/usr/sbin/guest_syncd
