# files/usr/lib/opkg/info/senao-systime.list
/etc/uci-defaults/systime
/usr/share/zoneinfo/timezone_list
/etc/init.d/systime
/usr/shc/timezonecheck
/sbin/httpsynctime.sh
/etc/config/systime
/sbin/checkInternet.sh
