# files/usr/lib/opkg/info/senao-snmp-utils.list
/usr/bin/snmpstatus
/usr/bin/snmptrap
/etc/init.d/snmpdtrap
/usr/bin/snmpset
/usr/bin/snmpget
/usr/bin/snmpwalk
/usr/sbin/snmpdtrap.sh
/usr/bin/snmptest
