# files/usr/lib/opkg/info/senao_capwap_c1.9.x.list
/etc/capwap/defaultWTPRoot.pem
/usr/sbin/upWtp.sh
/usr/bin/wtp
/usr/sbin/wtp_mem_check.sh
/usr/sbin/wsm.sh
/usr/sbin/capwap_sn_event.sh
/etc/capwap/config.wtp
/usr/sbin/test_ac.sh
/usr/sbin/conn_status.sh
/etc/init.d/apcontroller
/usr/sbin/portal.sh
/usr/sbin/check_code.sh
/etc/capwap/defaultWTP.pem
/etc/config/apcontroller
/etc/version_capwap
/etc/capwap/currentWTP.pem
/etc/capwap/currentWTPRoot.pem
/etc/capwap/currentWTPKey.pem
/etc/capwap/defaultWTPKey.pem
/etc/uci-defaults/z99capwap
