# files/usr/lib/opkg/info/senao-syskey.list
/etc/uci-defaults/z99syskey
/lib/syskey.sh
/etc/init.d/syskey
/usr/sbin/iface_mark.sh
/etc/config/syskey
/etc/init.d/iface_mark
