# files/usr/lib/opkg/info/senao-snmpd.list
/etc/init.d/snmpd
/etc/default/snmpd
/usr/shc/snmpd
/etc/config/snmpd
/usr/sbin/snmpd
/etc/snmp/snmpd.conf
/etc/uci-defaults/snmpd
