# files/usr/lib/opkg/info/senao-app-agent.list
/etc/init.d/app_agent
/etc/agent.pem
/usr/sbin/app_agentd
/etc/config/app_agent
