# files/usr/lib/opkg/info/senao-traffic-control.list
/etc/init.d/sntcd
/usr/bin/sntcd
/etc/assoc_notify.d/sntcd_notify.sh
/etc/config/sntcd
/etc/uci-defaults/sntcd
