# files/usr/lib/opkg/info/sn-websockets-server.list
/etc/init.d/swd
/usr/share/swd/swd.key.pem
/usr/bin/swd
/usr/share/swd/swd.pem
