# files/usr/lib/opkg/info/speedtest-master.list
/etc/init.d/speedtest
/sbin/speedtest-master
