# files/usr/lib/opkg/info/zzSENAO_VENDOR.list
/usr/shc/umount
/etc/init.d/setmac
/etc/init.d/cron
/usr/shc/factory_setting
/usr/shc/setmac
/usr/shc/setfastscan
/usr/shc/correctWifiConfig
/etc/init.d/led
/usr/shc/syslogd
/usr/shc/led
/etc/init.d/correctWifiConfig
/etc/init.d/syslogd
/usr/shc/boot
/etc/init.d/telnet
/etc/init.d/setfastscan
/usr/shc/telnet
/etc/init.d/umount
/etc/init.d/factory_setting
/usr/shc/ubootenvcfg
/etc/init.d/boot
/usr/shc/mount-nand-ubifs
/usr/shc/cron
